# supervisor.ini
[program:mattaudesse.com]
command=/home/mattaudesse/webapps/mattaudesse_com/mattaudesse.com serve
directory=/home/mattaudesse/webapps/mattaudesse_com
user=mattaudesse
autostart=true
autorestart=true
redirect_stderr=True
